# gw/etc/systemd/system/vpn-routing.service
[Unit]
Description=VPN routing init script

[Service]
ExecStart=/etc/openvpn/scripts/vpn-routing.sh

[Install]
WantedBy=multi-user.target

